# systemd/central-key-manager-listener.service.in
[Unit]
Description=Start the Central Key Manager Listener
Requires=dbus.service
After=central-key-manager.service

[Service]
User=@USER_NAME@
Group=@GROUP_NAME@
SmackProcessLabel=@SMACK_DOMAIN_NAME@
Type=simple
ExecStart=/usr/bin/key-manager-listener

[Install]
WantedBy=multi-user.target
