# systemd/key-manager.target
[Unit]
Description=Cert Key Manager sockets
DefaultDependencies=true

