# systemd/key-manager.target
[Unit]
Description=Central Key Manager sockets
DefaultDependencies=true

