# systemd/key-manager.service
[Unit]
Description=Start the Central Key Manager

[Service]
Type=notify
ExecStart=/usr/bin/key-manager

[Install]
WantedBy=multi-user.target
